# modules/utils/vcapture/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/output/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_RES)sx24 -ac +extension GLX +render -noreset
user=node
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/app/start-vnc.sh
user=node
environment=HOME="/home/node",USER="node",DISPLAY="%(ENV_DISPLAY)s",VNC_PASSWORD="%(ENV_VNC_PASSWORD)s"
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5930 --listen 8080
autostart=true
autorestart=true
priority=40
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pipewire]
command=/usr/bin/pipewire
environment=HOME="/home/node",XDG_RUNTIME_DIR="/run/user/1000"
user=node
autostart=true
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wireplumber]
command=/usr/bin/wireplumber
environment=HOME="/home/node",XDG_RUNTIME_DIR="/run/user/1000"
user=node
autostart=true
autorestart=true
priority=51
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pipewire-pulse]
command=/usr/bin/pipewire-pulse
environment=HOME="/home/node",XDG_RUNTIME_DIR="/run/user/1000"
user=node
autostart=true
autorestart=true
priority=52
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xfce4]
command=/usr/bin/xfce4-session
environment=HOME="/home/node",USER="node",DISPLAY="%(ENV_DISPLAY)s",XDG_RUNTIME_DIR="/run/user/1000"
user=node
autostart=true
autorestart=true
priority=70
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:audio-setup]
command=/app/setup-audio.sh
user=node
environment=HOME="/home/node",XDG_RUNTIME_DIR="/run/user/1000"
autostart=true
autorestart=false
startsecs=0
priority=55
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
